# Extra/clight.service
[Unit]
Description=C daemon to adjust screen brightness to match ambient brightness, as computed capturing frames from webcam

[Service]
Type=simple
ExecStart=/usr/bin/clight
Wants=clightd.service

[Install]
WantedBy=default.target
