# Extra/systemd/clight.service
[Unit]
Description=C daemon to adjust screen brightness to match ambient brightness, as computed capturing frames from webcam
Wants=clightd.service
After=display-manager.service clightd.service

[Service]
Type=simple
ExecStart=/usr/bin/clight
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
